# images/tdx-gpu/mkosi.extra/usr/lib/systemd/system/cvm-runtime-init.service
[Unit]
Description=CVM (TDX + NVIDIA CC) runtime initialisation
Documentation=https://github.com/Privasys/cvm-images
DefaultDependencies=no
# Run after /data is mounted (we read /data/nvidia-cc-bundle).
# Run before any platform service that needs the GPU.
After=local-fs.target data.mount
Requires=local-fs.target
Before=containerd.service manager.service caddy.service multi-user.target
# /dev/tdx_guest is the canonical TDX guest device on GCP/OVH/Linux >=6.2.
# /sys/firmware/tdx is NOT exposed by the in-tree TDX guest driver; do
# not use it as a condition (the unit will silently skip).
ConditionPathExists=/dev/tdx_guest
# Soft failure: don't block boot if NVIDIA setup fails.
# Operators can inspect /run/cvm-runtime-init.log.

[Service]
Type=oneshot
ExecStart=/usr/lib/cvm-runtime-init/cvm-runtime-init.sh
RemainAfterExit=yes
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
